# luci-app-overwall/root/etc/uci-defaults/luci-overwall
uci batch <<EOF
	delete ucitrack.@overwall[-1]
	add ucitrack overwall
	set ucitrack.@overwall[-1].init=overwall
	commit ucitrack
	delete firewall.overwall
	set firewall.overwall=include
	set firewall.overwall.type=script
	set firewall.overwall.path=/var/etc/overwall.include
	set firewall.overwall.reload=1
	commit firewall
        chmod +x /etc/init.d/overwall /usr/share/overwall/* >/dev/null 2>&1
	rm -rf /tmp/luci-modulecache/
	killall -HUP rpcd 2>/dev/null
	exit 0
EOF
